# cookbooks/zikraft/files/default/init/zikraft.conf
description "zikraft uwsgi service"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /home/www/moin/zikraft
exec /usr/local/bin/uwsgi -x /home/www/moin/zikraft/uwsgi.xml
respawn
