# cookbooks/inniti/templates/mysql/proxy_global_db.service.erb
# Generated by chef.
[Unit]
Description=<%= @description %>
After=network.target
[Service]
User=<%= @user || node['inniti']['user']['name'] %>
Type=forking
WorkingDirectory=<%= @working_directory %>
ExecStart=<%= @working_directory %><%= @execute_start %>
ExecStop=<%= @execute_stop %>
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
