# cookbooks/inniti/templates/flow_mes/flow_mes.service.erb
# Generated by chef.
[Unit]
Description=Fow_MES service
After=network.target
[Service]
User=<%= @user || node['inniti']['user']['name'] %>
Type=forking
WorkingDirectory=<%= @working_dir %>
ExecStart=<%= @working_dir %>/start.sh
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target